# scripts/linux/crontab
# Turn LED off at 22:00 every day
0 22 * * * /home/ccrsxx/homelab/scripts/linux/unifi-led.sh off

# Turn LED on at 07:00 every day
0 7 * * * /home/ccrsxx/homelab/scripts/linux/unifi-led.sh on
